# mkosi.extra/usr/lib/systemd/user/chezmoi-init.service
[Unit]
Description=Initializes Chezmoi from bundled dotfiles on first login
ConditionPathExists=!%h/.local/share/chezmoi

[Service]
ExecStart=/bin/bash -c 'cp -a /usr/share/harness/dotfiles %h/.local/share/chezmoi && chezmoi apply'
Type=oneshot

[Install]
WantedBy=graphical-session-pre.target
