# system_files/usr/lib/systemd/user/chezmoi-update.service
[Unit]
Description=Chezmoi Update

[Service]
ExecStart=chezmoi apply --no-tty --keep-going
Type=oneshot
