# system_files/etc/profile.d/hfetch.sh
alias neofetch="glorpfetch"
alias fastfetch="glorpfetch"
